# edge/raspberry-pi/systemd/komatsuna-agent.service
[Unit]
Description=Komatsuna AI Agent - Raspberry Pi observation service
After=network.target

[Service]
Type=oneshot
# git clone 先: ~/komatsuna-ai-agent（Pi ユーザー: kansei → .ssh-env 参照）
WorkingDirectory=/home/kansei/komatsuna-ai-agent/edge/raspberry-pi
EnvironmentFile=/home/kansei/komatsuna-ai-agent/edge/raspberry-pi/.env
ExecStart=/home/kansei/komatsuna-ai-agent/edge/raspberry-pi/.venv/bin/python -m app.main once

[Install]
WantedBy=multi-user.target
